# package/base-files/files/etc/hotplug.d/button/buttons_debug
#!/bin/sh
logger The button was $BUTTON and the action was $ACTION

